# ansible/roles/agent_auditor/files/agent-auditor.service
[Unit]
Description=Agent Auditor Service
After=network.target
Conflicts=apache2.service nginx.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agent-auditor
# Stop Apache or Nginx if running (they might be using port 80)
ExecStartPre=-/bin/bash -c "systemctl is-active apache2 && systemctl stop apache2 || true"
ExecStartPre=-/bin/bash -c "systemctl is-active nginx && systemctl stop nginx || true"
# Set TARGET_HOST to the machine's own IP address (script is extracted from the binary at runtime)
ExecStartPre=/bin/bash -c 'cd /opt/agent-auditor && ./scripts/set_target_host.sh /opt/agent-auditor/.env'
ExecStart=/opt/agent-auditor/agent-auditor
Restart=on-failure
RestartSec=5
Environment=PORT=8084
Environment=PROXY_PORT=80
Environment=AEGONG_DEV_MODE=0

# Load environment variables from .env file if it exists
EnvironmentFile=-/opt/agent-auditor/.env

[Install]
WantedBy=multi-user.target